# _/udev/10-autorandr.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/admin/.Xauthority", RUN+="autorandr --change --default common"
